# fs/etc/services.d/download/run
#!/usr/bin/env python3

from os import execl
from pathlib import Path
from sys import executable

execl(
    executable,
    executable,
    Path().cwd() / "download.py",
)
