# fs/etc/services.d/broadwayd/run
#!/bin/bash

/command/s6-notifyoncheck /bin/broadwayd --port $BROADWAY_PORT
